# rootfs/etc/services.d/automator-server/run
#!/usr/bin/with-contenv bash
#python /automator/server.py

